# systemd/user/xapblr@.service
[Unit]
Description=Locally index a Tumblr blog

[Service]
Type=oneshot
ExecStart=xapblr index %i
