# systemd/user/xapblr-daily@.timer
[Unit]
Description=Locally index a Tumblr blog, daily.

[Timer]
OnCalendar=daily
Persistent=true
Unit=xapblr@%i.service
