# systemd/user/xapblr-hourly@.timer
[Unit]
Description=Locally index a Tumblr blog, hourly.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
